# cyberpanel/seconddns.service
[Unit]
Description=SecondDNS domain watcher for CyberPanel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/seconddns watch
Restart=always
RestartSec=10
StandardOutput=append:/var/log/seconddns.log
StandardError=append:/var/log/seconddns.log

[Install]
WantedBy=multi-user.target
